# systemd/auto-scrapy-runtime.service
[Unit]
Description=Run the auto-scrapy pipeline once
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/247-data-collector
Environment=AUTO_SCRAPY_SOURCES_CONFIG_PATH=config/sources/target_smoke_sources.toml
ExecStart=/home/xuan/xuanStudy/env-utils/miniforge3/envs/247-data-collector/bin/uv run python -m app.runtime

# Optional: Ensure Python logs appear immediately in journalctl
Environment=PYTHONUNBUFFERED=1
