# systemd/backup/auto-scrapy-runtime.service
[Unit]
Description=Run the auto-scrapy pipeline once
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/247-data-collector
Environment=AUTO_SCRAPY_SOURCES_CONFIG_PATH=config/sources/target_smoke_sources.toml
ExecStart=/usr/bin/env sh -lc 'uv run python -m app.runtime'
    